# maxwell-3306.conf
description "Maxwell binlog-to-kafka bridge"

start on filesystem and net-device-up IFACE=lo
respawn limit 15 5

env CONF="/etc/mysql/maxwell-3306.conf"

pre-start script
  [ -r $CONF ] || {
     echo "I can't find or read $CONF">&2
     exit 1
  }
end script

console log
chdir /usr/local/bin/maxwell/maxwell-1.8.1

script
  echo $$ > /var/run/maxwell-3306.pid
  exec /usr/local/bin/maxwell/maxwell-1.8.1/bin/maxwell --config /etc/mysql/maxwell-3306.conf
end script

post-stop script
  rm -f /var/run/maxwell-3306.pid
end script
